# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn Server for SITENAME"

start on net-device-up
stop on shutdown

respawn

chdir /home/elspeth/sites/SITENAME/source

exec /home/elspeth/env1/bin/gunicorn --bind unix:/tmp/SITENAME.socket superlists.wsgi:application